# next/3-settings-lxde/debian/paintitblack-lxde-livesettings.service
[Unit]
Description=LXDE cheatcodes in live media
ConditionPathExists=/usr/share/paintitblack-lxde-livesettings/
After=paintitblack-common-livesettings.service

[Service]
Type=oneshot
ExecStart=/usr/share/paintitblack-lxde-livesettings/lxde-cheatcodes start

[Install]
WantedBy=graphical.target
